# gebo.apps.parent/gebo.ai.app/src/packaging/deb/resources/postinst
#!/bin/bash
set -e

APP_USER="gebo-ai"
APP_GROUP="gebo-ai"
APP_HOME="/home/gebo-ai"
GEBO_HOME="$APP_HOME/home"
GEBO_WORK="/var/gebo-ai/data"
SERVICE_NAME="gebo-ai.service"
APP_DIR="$APP_HOME/gebo-ai"
SERVICE_SRC="$APP_DIR/$SERVICE_NAME"
SERVICE_DST="/etc/systemd/system/$SERVICE_NAME"

# Create user and service
if ! id -u "$APP_USER" >/dev/null 2>&1; then
    useradd -r -m -d "$APP_HOME" -s /bin/bash "$APP_USER"
fi

# Create folders
mkdir -p "$GEBO_HOME" "$GEBO_WORK"
chown -R "$APP_USER:$APP_GROUP" "$GEBO_HOME" "$GEBO_WORK"

# Install systemd unit
if [ -f "$SERVICE_SRC" ]; then
    cp "$SERVICE_SRC" "$SERVICE_DST"
    chmod 744 "$SERVICE_DST"
fi

# Reload systemd and enables the service
if command -v systemctl >/dev/null  then
    systemctl daemon-reload || true
    systemctl enable "$SERVICE_NAME" || true
    # opzionale: avvio immediato
    systemctl restart "$SERVICE_NAME" || true
fi

exit 0
